# data/user-identification-manager.service.in
[Unit]
Description=User Identification Manager

[Service]
Type=dbus
BusName=com.luxoft.UserIdentificationManager
ExecStart=@bindir@/user-identification-manager
